# contrib/turnup-ui.service
[Unit]
Description=Turn Up web UI
Documentation=https://github.com/sean351/turn-up-arch
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/turnup-ui
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=turnup-ui

[Install]
WantedBy=default.target
